# sia-receiver.service
[Unit]
Description=SIA Protocol Receiver Daemon
Documentation=https://github.com/delgod/ajax-imou-bridge
After=network.target network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=nobody
Group=nogroup
WorkingDirectory=/tmp
EnvironmentFile=/etc/sia-receiver.conf
ExecStart=/usr/local/bin/sia_receiver.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sia-receiver

# Resource limits
LimitNOFILE=65536
MemoryHigh=256M
MemoryMax=512M

[Install]
WantedBy=multi-user.target 
